# packages/hostha/hostha-0.0.0.tar.gz/hostha-0.0.0/package/systemd/hostha-api.service
[Unit]
Description=Hostha API Service
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/hostha-api --config-file /etc/hostha/hostha.conf --log-file /var/log/hostha/api.log
Restart=always

[Install]
WantedBy=multi-user.target
